# utils/99-petitboot.rules
# tell petitboot when we see new block devices ...
SUBSYSTEM=="block",RUN+="/usr/sbin/udev-helper"
